# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:fuseki]
command=/bin/bash /opt/fuseki/start-fuseki.sh


